# deployment/systemd/secrets-check.service
# /etc/systemd/system/secrets-check.service
# Stellt sicher dass /run/secrets/ befüllt ist bevor Docker startet.
# Nach Reboot: Docker blocked bis CI/CD die Secrets erneut deployed.
#
# Installation auf odoo-prod (46.225.127.211):
#   scp deployment/systemd/secrets-check.service root@46.225.127.211:/etc/systemd/system/
#   ssh root@46.225.127.211 "systemctl daemon-reload && systemctl enable secrets-check"
#
# Philosophie (ADR-045): kein Secret = kein Service.
# Verhindert Silent-Startup mit fehlenden Credentials nach Server-Reboot.
# /run/ ist tmpfs — wird bei Reboot geleert. Neustart erfordert CI/CD-Re-Deploy.

[Unit]
Description=Verify /run/secrets/ populated before Docker starts (ADR-045)
Before=docker.service
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c '\
  REQUIRED="postgres_password postgres_db postgres_user odoo_admin_passwd odoo_domain acme_email"; \
  MISSING=0; \
  for f in $REQUIRED; do \
    if [ ! -s "/run/secrets/$f" ]; then \
      echo "FATAL: Missing /run/secrets/$f" >&2; \
      MISSING=$((MISSING + 1)); \
    fi; \
  done; \
  if [ $MISSING -gt 0 ]; then \
    echo "odoo-hub cannot start: $MISSING secret(s) missing." >&2; \
    echo "Re-deploy via GitHub Actions to populate /run/secrets/." >&2; \
    exit 1; \
  fi; \
  echo "odoo-hub secrets verified ($(ls /run/secrets/ | wc -l) files in /run/secrets/)."; \
'

[Install]
WantedBy=multi-user.target
